# Linux/Logs/crontab

*/5 * * * * php /home/boc_a/Server20170304/Linux/Logs/importLog2DB.php >/tmp/crontab.log
